# .config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland Wayland Session
Wants=waybar.service
Wants=swaync.service
Wants=hypridle.service
