# .config/systemd/user/quickshell-lock-monitor.timer
[Unit]
Description=Timer to check quickshell lock socket health
After=quickshell.service

[Timer]
# Check every 5 minutes
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
