# .config/systemd/user/waybar-reload.service
[Unit]
Description=Reload Waybar configuration

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user reload-or-restart waybar.service