# .config/systemd/user/waybar-reload.path
[Unit]
Description=Watch Waybar config directory for changes
Requires=waybar.service

[Path]
PathModified=/home/bennyp/.config/waybar
Unit=waybar-reload.service

[Install]
WantedBy=waybar.service