# .config/systemd/user/noctalia.service.d/override.conf
[Service]
Environment=QT_SCALE_FACTOR=2
Environment=QT_QPA_PLATFORMTHEME=qt6ct

[Unit]
Wants=hypridle.service
