# .config/systemd/user/quickshell-polkit-agent.service.d/override.conf
[Service]
Environment=QT_LOGGING_RULES="*.debug=true"
